# vibora-dev.service
[Unit]
Description=Vibora Development Server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/projects/vibora
Environment=VIBORA_DEVELOPER=1
Environment=NODE_ENV=production
# Note: Build is handled by the /api/config/restart endpoint before triggering restart
# This ensures we only restart if build succeeds
ExecStartPre=/bin/bash -c 'cd %h/projects/vibora && bun run drizzle-kit push'
ExecStart=/bin/bash -c 'cd %h/projects/vibora && bun server/index.ts'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
